# systemd/ptt-dictation.service
[Unit]
Description=Push-to-Talk Dictation (ElevenLabs Scribe)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/SynologyDrive/Dev/11labs/ptt-dictation
ExecStart=/usr/bin/python3 -m ptt_dictation
Restart=on-failure
# Optionally provide the key here instead of config.toml:
# Environment=ELEVENLABS_API_KEY=sk-...

[Install]
WantedBy=graphical-session.target
